# supervisor.conf
[program:EthermineTelegramBots]
command=/opt/EthermineTelegramBots/.venv/bin/python eth.py
directory=/opt/EthermineTelegramBots
user=nobody
autorestart=true
redirect_stderr=true